# agent/linux/systemd/tuxcmdb-agent.service
[Unit]
Description=TuxCMDB Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/tuxcmdb-agent --once

[Install]
WantedBy=multi-user.target
